# xbindkeys/xbindkeys.conf
# xbindkeys - set custom key-bindings
#
# This task is run on startup to set the key-bindings from ~/.xbindkeysrc

description     "set custom key-bindings"

start on startup

exec xbindkeys

